# py_backend/service/iris-backend.service
# IRIS Command Backend Systemd Service
#
# Installation:
#   sudo cp iris-backend.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable iris-backend
#   sudo systemctl start iris-backend
#
# Commands:
#   sudo systemctl start iris-backend
#   sudo systemctl stop iris-backend
#   sudo systemctl restart iris-backend
#   sudo systemctl status iris-backend
#   journalctl -u iris-backend -f   # Follow logs

[Unit]
Description=IRIS Command Video Analytics Backend
Wants=network-online.target iris-mediamtx.service
After=network-online.target iris-mediamtx.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/Desktop/iris_command/py_backend

# Critical: Unbuffered Python output for streaming to work
Environment=PYTHONUNBUFFERED=1
Environment=OPENCV_FFMPEG_CAPTURE_OPTIONS=rtsp_transport;tcp|buffer_size;262144|analyzeduration;50000|probesize;50000|fflags;nobuffer|flags;low_delay

# Use the wrapper to ensure timestamped logs in /home/ubuntu/Desktop/iris_command/logs/YYYY-MM-DD/
ExecStart=/home/ubuntu/Desktop/iris_command/py_backend/service/run_backend.sh

# Restart policy
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=6

# Logging
StandardOutput=null
StandardError=null

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
