# py_backend/service/iris-backend.service
# IRIS Command Backend Systemd Service
#
# Installation:
#   sudo cp iris-backend.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable iris-backend
#   sudo systemctl start iris-backend
#
# Commands:
#   sudo systemctl start iris-backend
#   sudo systemctl stop iris-backend
#   sudo systemctl restart iris-backend
#   sudo systemctl status iris-backend
#   journalctl -u iris-backend -f   # Follow logs

[Unit]
Description=IRIS Command Video Analytics Backend
Requires=iris-mediamtx.service
Wants=network-online.target
After=network-online.target iris-mediamtx.service
PartOf=iris-mediamtx.service
StartLimitIntervalSec=60
StartLimitBurst=6

[Service]
Type=simple
User=ubuntu
Group=ubuntu
SupplementaryGroups=video render
WorkingDirectory=/home/ubuntu/Desktop/iris_command/py_backend

# Critical: Unbuffered Python output for streaming to work
Environment=PYTHONUNBUFFERED=1
Environment=OPENCV_FFMPEG_CAPTURE_OPTIONS=rtsp_transport;tcp|buffer_size;262144|analyzeduration;50000|probesize;50000|fflags;nobuffer|flags;low_delay
EnvironmentFile=-/etc/default/iris-command
ExecStartPre=/usr/bin/test -x /home/ubuntu/Desktop/iris_command/py_backend/.venv/bin/python
ExecStartPre=/usr/bin/test -f /home/ubuntu/Desktop/iris_command/py_backend/app.py

# Use the wrapper to ensure timestamped logs in /home/ubuntu/Desktop/iris_command/logs/YYYY-MM-DD/
ExecStart=/home/ubuntu/Desktop/iris_command/py_backend/service/run_backend.sh

# Restart policy
Restart=always
RestartSec=3
TimeoutStopSec=45
KillSignal=SIGINT
KillMode=mixed

# Logging
StandardOutput=journal
StandardError=journal

# Resource limits
LimitNOFILE=65535
TasksMax=infinity

[Install]
WantedBy=multi-user.target
